# jira-daily-planner.cron
0 1 * * * python /usr/src/app/jira-daily-planner.py > /proc/1/fd/1 2>/proc/1/fd/2
